# owncloud-db/owncloud-db-backup.cron
#!/bin/bash
/root/bin/owncloud-db-backup.sh 2>&1
